# ubuntu-gnome-15.10/displaylink/80-displaylink-powersave.rules
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/displaylink-powersave.sh true"
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/local/bin/displaylink-powersave.sh false"
